# scie2100.conf
description "uWSGI server instance configured to serve myproject"

start on runlevel [2345]
stop on runlevel [!2345]

setuid gabe
setgid gabe

env PATH=/home/gabe/quiz
chdir /home/gabe/quiz
exec make
